# scripts/com.nanoclaw.news-summary.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nanoclaw.news-summary</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/env</string>
        <string>npx</string>
        <string>tsx</string>
        <string>scripts/news-summary.ts</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/huanyu.guo/self/nanoclaw</string>
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>20</integer>
    </dict>
    <key>EnvironmentVariables</key>
    <dict>
        <key>TZ</key>
        <string>Asia/Shanghai</string>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        <key>HOME</key>
        <string>/Users/huanyu.guo</string>
    </dict>
    <key>StandardOutPath</key>
    <string>/Users/huanyu.guo/self/nanoclaw/logs/news-summary-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/huanyu.guo/self/nanoclaw/logs/news-summary-stderr.log</string>
</dict>
</plist>
